# prototype/whisper-p2t.service
# Template — setup.sh generates the actual service file with correct paths.
# Manual install: update paths and env vars below, then:
#   sudo cp whisper-p2t.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now whisper-p2t

[Unit]
Description=Push-to-Talk Dictation Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/path/to/dictate/venv/bin/python /path/to/dictate/push2talk.py
Environment=SUDO_USER=YOUR_USER
Environment=SUDO_UID=YOUR_UID
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/YOUR_UID
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
